# supervisord.conf
[program:node]
environment = NODE_ENV=production
command = node index.js
directory = /home/dotcloud/current